# conf/supervisord.conf
[supervisord]
nodaemon=true
childlogdir=/var/log
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:nginx]
command = /usr/sbin/nginx
user = root
autostart = true
log_stdout=true
log_stderr=true
logfile=/var/log/nginx.log

[program:php-fpm]
command = /usr/sbin/php-fpm -FR
user = root
autostart = true
log_stdout=true
log_stderr=true
logfile=/var/log/php-fpm.log

[program:crond]
command=/usr/sbin/crond -n
log_stdout=true
log_stderr=true
logfile=/var/log/crond.log

[program:crontab]
command=/usr/bin/crontab
log_stdout=true
log_stderr=true
logfile=/var/log/crontab.log

[program:sshd]
command=/usr/sbin/sshd -D
log_stdout=true
log_stderr=true
logfile=/var/log/sshd.log
